# dictd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dictd --listen-to 127.0.0.1 --locale en_US.UTF-8 -- -s
