# dictd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dictd --listen-to localhost --locale en_US.UTF-8 -- -s
